# 04/keyboard.rules
ACTION=="add", SUBSYSTEM=="input", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceSubClass}=="01", ATTRS{bInterfaceProtocol}=="01", RUN+="/usr/sbin/modprobe /home/jkate/04/main.ko"
